# grob.container
# grob.container
# Quadlet systemd container unit for rootless Podman
# Install: cp grob.container ~/.config/containers/systemd/
# Enable: systemctl --user daemon-reload && systemctl --user enable --now grob

[Unit]
Description=Grob LLM Routing Proxy
Documentation=https://github.com/gelwood/grob
Wants=network-online.target
After=network-online.target

[Container]
# Image configuration
Image=grob:latest
Pull=never

# Networking
PublishPort=8080:8080
PublishPort=9090:9090

# User/Group (rootless)
User=1000
Group=1000

# Security options
NoNewPrivileges=true
DropCapabilities=ALL
AddCapability=NET_BIND_SERVICE
SecurityLabelType=container_runtime_t
ReadOnly=true

# Environment
Environment=RUST_LOG=info
Environment=GROB_CONFIG_PATH=/etc/grob/config.toml
Environment=GROB_DATA_DIR=/var/lib/grob

# Volumes (host:container)
Volume=%h/.grob/config:/etc/grob:ro,Z
Volume=%h/.grob/data:/var/lib/grob:Z
Volume=grob-tmp:/tmp

# Health check
HealthCmd=/grob status --format json
HealthInterval=30s
HealthTimeout=5s
HealthStartPeriod=5s
HealthRetries=3

# Resource limits
Memory=512m
MemorySwap=512m
CPUQuota=100%
PidsLimit=1024

[Service]
Type=notify
NotifyAccess=all
Restart=on-failure
RestartSec=5
TimeoutStartSec=60
TimeoutStopSec=30

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=grob

[Install]
WantedBy=default.target
