# deploy/grob.container
# grob.container
# Quadlet systemd container unit for rootless Podman
# Install: cp deploy/grob.container deploy/grob.volume ~/.config/containers/systemd/
# Enable: systemctl --user daemon-reload && systemctl --user enable --now grob

[Unit]
Description=Grob LLM Routing Proxy
Documentation=https://github.com/azerozero/grob
Wants=network-online.target
After=network-online.target

[Container]
# Image configuration
Image=grob:latest
Pull=never

# Networking
PublishPort=8080:8080
PublishPort=9090:9090

# User/Group matches the scratch image user (65534 = nobody). The data bind
# mount below uses :U so rootless Podman makes it writable for this UID.
User=65534
Group=65534

# Security options
NoNewPrivileges=true
DropCapabilities=ALL
SecurityLabelType=container_runtime_t
ReadOnly=true

# Environment
Environment=RUST_LOG=info
Environment=GROB_CONFIG=/etc/grob/config.toml
Environment=GROB_HOME=/var/lib/grob

# Volumes (host:container)
Volume=%h/.grob/config:/etc/grob:ro,Z
Volume=%h/.grob/data:/var/lib/grob:Z,U
Volume=grob-tmp:/tmp

# Health check
HealthCmd=/grob status
HealthInterval=30s
HealthTimeout=5s
HealthStartPeriod=5s
HealthRetries=3

# Resource limits
Memory=512m
MemorySwap=512m
CPUQuota=100%
PidsLimit=1024

[Service]
# NOTE: Type=exec over Type=simple — systemd waits for execve() to succeed
# before considering the unit started, surfacing exec-time failures (missing
# image, invalid args) as real start failures instead of immediate-exit loops.
# Type=notify is wrong here: grob does not call sd_notify, so systemd would
# wait for a notification that never arrives. Readiness is signaled by the
# Quadlet HealthCmd above.
Type=exec
Restart=on-failure
RestartSec=5
TimeoutStartSec=60
TimeoutStopSec=30

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=grob

[Install]
WantedBy=default.target
